# deploy/com.knowledgereduce.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.knowledgereduce.watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/knowledgereduce</string>
        <string>watch-daemon</string>
        <key>--dir</key>
        <string>/Users/nelslindahl/knowledgereduce_watch</string>
        <key>--store</key>
        <string>/Users/nelslindahl/knowledgereduce_store</string>
        <key>--db-log</key>
        <string>/Users/nelslindahl/knowledgereduce_watcher.db</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nelslindahl/Library/Logs/knowledgereduce_watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nelslindahl/Library/Logs/knowledgereduce_watcher.log</string>
</dict>
</plist>
